# openrc/init.d/squashmount
#!/sbin/runscript
# (C) Martin V\"ath <martin@mvath.de>

depend() {
	need localmount
}

description='mount/umount all squashmount configured mountpoints'

start() {
	ebegin 'mounting all squashmount configure mountpoints'
	squashmount mount
	eend ${?}
	:
}

stop() {
	ebegin 'umounting all squashmount configure mountpoints'
	squashmount umount
	eend ${?}
	:
}
