# openrc/init.d/squashmount
#!/sbin/openrc-run
# (C) Martin V\"ath <martin at mvath.de>
# SPDX-License-Identifier: BSD-3-Clause

depend() {
	need localmount tmpfiles.setup
}

description='start/stop all squashmount configured mount points'

start() {
	ebegin 'mounting with squashmount'
	squashmount start
	eend $?
}

stop() {
	ebegin 'umounting/stopping squashmount'
	squashmount -f --lsof=0 --lsof-ro=0 stop
	eend $?
}
