# systemd/system/squashmount.service
[Unit]
Description=mount/umount all squashmount configured mountpoints
After=local-fs.target systemd-tmpfiles-setup.service

[Service]
Type=oneshot
RemainAfterExit=true
ExecStart=/usr/bin/squashmount start
ExecStop=/usr/bin/squashmount -f stop

[Install]
WantedBy=multi-user.target
