# etc/systemd/system/squashmount.service.d/timeout.conf
[Service]
# Replace the value by a generous upper bound for the maximally
# expected compression time on your system:
TimeoutStopSec=15m
